# files/timers/btrfs-balance-root.service
[Service]
Type=oneshot
ExecStart=/bin/btrfs balance start -dusage=50 -musage=50 /
