# files/timers/btrfs-balance-root.service
[Service]
Type=oneshot
ExecStart=/bin/btrfs balance start -dusage=75 -musage=75 /
