# files/timers/cache-updates.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/cache-updates
