# files/timers/fstrim-ext.service
[Unit]
ConditionPathIsMountPoint=/srv/samba

[Service]
Type=oneshot
ExecStart=/bin/fstrim -v /srv/samba
