# files/timers/fstrim-root.service
[Service]
Type=oneshot
ExecStart=/bin/fstrim -v /
