# files/timers/hwclock-systohc.service
[Service]
Type=oneshot
ExecStart=/usr/sbin/hwclock --systohc
