# files/timers/keep-disk-on.timer
[Timer]
OnCalendar=*:0/2
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
