# files/timers/keep-disk-on.service
[Unit]
ConditionPathExists=/srv/samba/.keepalive

[Service]
Type=oneshot
ExecStart=/bin/touch /srv/samba/.keepalive
