# files/timers/munin-cron.service
[Service]
Type=oneshot
User=munin
ExecStart=/usr/bin/munin-cron
