# files/timers/pacman-fy.service
[Service]
Type=oneshot
ExecStart=/bin/pacman -Fy
