# files/timers/pacman-fy.service
[Unit]
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecCondition=/usr/local/bin/internet-isnt-crappy
ExecStart=/bin/pacman -Fy
