# files/timers/restart-atop-midnight.service
[Service]
Type=oneshot
ExecStart=/usr/sbin/systemctl restart atop-oneday
