# files/timers/unbound-root-hints-update.service
[Service]
Type=oneshot
ExecStart=/bin/curl -o /etc/unbound/root.hints https://www.internic.net/domain/named.cache
