# files/timers/unbound-root-hints-update.timer
[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=1d

[Install]
WantedBy=timers.target
