# files/timers/vacuum-journals.service
[Service]
Type=oneshot
ExecStart=/bin/journalctl --rotate
ExecStart=/bin/journalctl --vacuum-time=30d
