# files/timers/vacuum-journals.timer
[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1d

[Install]
WantedBy=timers.target
