# files/ansible-update.timer
[Timer]
OnCalendar=daily
RandomizedDelaySec=6h

[Install]
WantedBy=multi-user.target
