# files/ansible-update.timer
[Timer]
OnCalendar=*:0/15
Unit=ansible-update.service

[Install]
WantedBy=multi-user.target
