# files/ansible-update.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/ansible-update

[Install]
WantedBy=multi-user.target
