# files/enable-zswap-or-zram.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/enable-zswap-or-zram

[Install]
WantedBy=multi-user.target
