# files/l4d2-restart.service
[Unit]
Description=Left 4 Dead 2 periodic restart

[Service]
ExecStart=/usr/bin/systemctl restart l4d2
