# files/night-sleep.service
[Unit]
Description=Put computer to sleep

[Service]
ExecStart=/bin/sh -c 'printf mem > /sys/power/state && sleep 20'
