# files/night-sleep.service
[Unit]
Description=Sleep for 10 hours

[Service]
ExecStart=/usr/bin/rtcwake -s 36000 -m mem
