# files/udev/90-wifi-powersave.rules
ACTION=="add|change", SUBSYSTEM=="net", KERNEL=="wl*", RUN+="/usr/bin/iw dev %k set power_save on"
